# keepalived/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true

[program:keepalived]
command=keepalived -n --log-console
autostart=true
